# ubuntu/odoo/supervisord.conf
[program:odoo]
command=/opt/odoo/start.sh
user=root
